# salt/oscarcommerce_base/source/etc/supervisor.conf
{% set project_name = salt['pillar.get']('project_name') %}
{% set admin_user = salt['pillar.get']('users:admin_user') %}

[program:{{ project_name }}]
command = /opt/{{ project_name }}/gunicorn.sh
user = {{ admin_user }}
stdout_logfile = /opt/{{ project_name }}/log/supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
